# containers/twitch/files/root/etc/services.d/encoder/run
#!/usr/bin/with-contenv bash

cd /app || exit

[[ ${START_ENC:-1} -eq 0 ]] && exit

exec s6-setuidgid abc python launcher.py $GLOBAL_ARGS encoder $ENC_ARGS
