# promscale/promscale-maintaince.timer
[Unit]
Description=Promscale maintenance timer

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
