# twitch/files/root/etc/services.d/encoder/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

[[ -S ENC_LISTEN_ADDRESS ]] && rm -fv "$ENC_LISTEN_ADDRESS"
